# templates/borgmatic/borgmatic.timer.j2
[Unit]
Description=Run borgmatic backup for {{ service_name }}

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
